# client/packaging/macos/com.dtsys.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dtsys.agent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/dtsys-agent</string>
        <string>--config</string>
        <string>/Library/Application Support/dtsys/agent.toml</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/dtsys-agent.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/dtsys-agent.log</string>
</dict>
</plist>
